# mnt-NextCLoud.automount
[Unit]
Description=Automount router USB share for Nextcloud

[Automount]
Where=/mnt/NextCloud

[Install]
WantedBy=multi-user.target
